# home/config/systemd/user/ical-sync.timer
[Unit]
Description=Keep local calendars up to date

[Timer]
OnActiveSec=1m
OnUnitInactiveSec=5m

[Install]
WantedBy=timers.target
