# home/config/systemd/user/imap-sync.service
[Unit]
Description=Keep local mailboxes up to date
After=network.target network-online.target dbus.socket

[Service]
Type=oneshot
ExecStart=/usr/bin/env offlineimap
TimeoutStartSec=10min

[Install]
WantedBy=default.target
