# home/config/systemd/user/imap-sync-quick.timer
[Unit]
Description=Keep local mailboxes up to date

[Timer]
OnActiveSec=1m
OnUnitInactiveSec=5m

[Install]
WantedBy=timers.target
