# home/config/systemd/user/imap-sync-full.timer
[Unit]
Description=Keep local mailboxes up to date

[Timer]
OnActiveSec=5m
OnUnitInactiveSec=60m

[Install]
WantedBy=timers.target
