# home/config/systemd/user/imap-sync.timer
[Unit]
Description=Keep local mailboxes up to date

[Timer]
OnUnitInactiveSec=5m

[Install]
WantedBy=timers.target
