# system/etc/udev/rules.d/10-monitor.rules
# on monitor connect/disconnect, configure xrandr & polybar
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/will/bin/runx /home/will/bin/msu -s 1"
