# system/etc/udev/rules.d/99-wifi-nfs.rules
# on wifi changes, disable or enable NFS mounts
ACTION=="add|change", SUBSYSTEM=="net", ENV{DEVTYPE}=="wlan", RUN+="/home/will/bin/runx /usr/local/bin/wifi-nfs.sh"
