# system/etc/systemd/system/mnt-mother-data.automount
[Unit]
Description=Automount the mother NAS at home

[Automount]
Where=/mnt/mother/data

[Install]
WantedBy=multi-user.target
