# system/etc/systemd/system/pac-gen-mirrorlist.timer
[Unit]
Description=Update the pacman mirrorlist

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
